# ressources/systemd/user/encryptsync-queue.path
[Unit]
Description=Watch EncryptSync session queue

[Path]
MakeDirectory=true
DirectoryMode=0700
PathExistsGlob=%t/encryptsync/open-*
PathExistsGlob=%t/encryptsync/close-*
Unit=encryptsync-dispatch.service

[Install]
WantedBy=default.target