# ressources/systemd/user/encryptsync.service
[Unit]
Description=EncryptedSync - clear plaintext on logout/shutdown
After=default.target
Conflicts=exit.target
Before=exit.target

[Service]
Type=oneshot
WorkingDirectory=/usr/lib/encryptsync
RemainAfterExit=yes
ExecStart=/usr/bin/true
ExecStop=/usr/bin/python3 /usr/lib/encryptsync/main.py
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/bin:/bin
StandardOutput=journal
StandardError=journal
TimeoutStopSec=180

[Install]
WantedBy=default.target