# ressources/systemd/user/encryptsync-clear@.service
[Unit]
Description=EncryptedSync - clear plaintext on logout per-session (%i)

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/true
ExecStop=/usr/bin/encryptsyncctl clear --yes
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
TimeoutStopSec=180