# ressources/systemd/user/encryptsync-dispatch.service
[Unit]
Description=Process EncryptSync session queue

[Service]
Type=oneshot
ExecStart=/usr/lib/encryptsync/scripts/dispatch.sh