# service/systemd/encrypted-sync-daemon.service
[Unit]
Description=EncryptedSync daemon (main watcher)
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/project
ExecStart=/usr/bin/env python3 /path/to/project/main.py
Restart=always
RestartSec=3
Environment="PATH=/path/to/venv/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
