# ressources/systemd/system/encryptsync-clear.service
[Unit]
Description=EncryptedSync - clear plaintext on shutdown
DefaultDependencies=no
Before=shutdown.target
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/bin/encryptsyncctl clear --yes
TimeoutStopSec=180
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/bin:/bin

[Install]
WantedBy=multi-user.target